# year3/SWE/ngUML.component.backend/contrib/supervisord.conf
[program:nguml]
command=./nguml/manage.py serve
stderr_logfile = /var/log/supervisord/nguml-stderr.log
stdout_logfile = /var/log/supervisord/nguml-stdout.log